# server/notify-multiplex-server.in
#!/sbin/runscript
# This script starts and stops the notify-multiplex server

prefix=@prefix@
exec_prefix=@exec_prefix@
sysconfdir=@sysconfdir@
sbindir=@sbindir@

extra_started_commands=""

depend() {
	need net
	need logger
	after iptables
}

start() {
	ebegin "Starting notify-multiplex-server"
	start-stop-daemon -b -m --pidfile @localstatedir@/notify-multiplex-server.pid \
		-- @sbindir@/notify-multiplexd \
		@sysconfdir@/notify-multiplexer/notify-multiplexer.conf
	eend $? "Failed to start notify-multiplex-server"
}

stop() {
	ebegin "Stopping notify-multiplex-server"
	start-stop-daemon --stop -p @localstatedir@/notify-multiplex-server.pid
	eend $? "Failed to stop notify-multiplex-server"
}
